# !PtokaX 0.5.2.2 - CC LEDE - feeds patch makefile cfg/автозапуск/ptokax/ptokax
#!/bin/sh /etc/rc.common

START=98

start() {
 PtokaX -d -c /root/.PtokaX
}

stop() {
 killall PtokaX > /dev/null
}

restart() {
 killall PtokaX > /dev/null
 sleep 1
 PtokaX -d -c /root/.PtokaX
}
